# deploy/systemd/stock-analyze-polling.service
[Unit]
Description=stock-analyze polling
After=network.target

[Service]
Type=simple
WorkingDirectory=/app/stock-analyze
EnvironmentFile=-/app/stock-analyze/.env
ExecStart=/app/stock-analyze/.venv/bin/python -m app
Restart=always
RestartSec=5
StandardOutput=append:/app/stock-analyze/logs/polling.log
StandardError=append:/app/stock-analyze/logs/polling.log

[Install]
WantedBy=multi-user.target
